# deploy/hr-queue-worker.service
[Unit]
Description=HR Robot Queue Worker
After=network.target mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/hr-brb
ExecStart=/usr/bin/php artisan queue:work --sleep=3 --tries=3 --max-time=3600
Restart=always
RestartSec=5

# Логи
StandardOutput=append:/var/log/hr-robot/queue-worker.log
StandardError=append:/var/log/hr-robot/queue-worker-error.log

[Install]
WantedBy=multi-user.target
